# save_wiki.cron
00 19 * * * root /usr/bin/python /bin/save_wiki.py